# crontab
*/5     *       *       *       *       /root/startvnc.sh
